# output-hass/files/etc/uci-defaults/output-hass
#!/bin/sh
# Karl Palsson, <karlp@tweak.net.au> Nov 2021
touch /etc/config/output-hass
exit 0
